# rootfs/etc/services.d/bootstrap-script/run
#!/usr/bin/with-contenv sh
apache2ctl -D FOREGROUND